# 02_Operation_of_Running_Systems/A07/backup.service
[Unit]
Description="Back up the etc directory"

[Service]
ExecStart=tar -czf /root/etc.tar.gz /etc